# clinic-bot.service
[Unit]
Description=ROYAL Clinic Telegram Bot
After=network.target

[Service]
Type=simple
User=clinicbot
WorkingDirectory=/opt/clinic-bot
Environment="PATH=/opt/clinic-bot/venv/bin"
ExecStart=/opt/clinic-bot/venv/bin/python /opt/clinic-bot/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=clinic-bot

# Ограничения ресурсов
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
